# docker/s6-overlay/s6-rc.d/svc-grackle/run
#!/command/with-contenv bash

# Source all s6 container_environment vars explicitly — with-contenv should
# inject them, but gosu's exec can lose them on some platforms.
for f in /run/s6/container_environment/*; do
  [ -r "$f" ] || continue
  varname="$(basename "$f")"
  varvalue="$(<"$f")"
  export "$varname"="$varvalue"
done

exec gosu node grackle serve --allow-network
