# docker/s6-overlay/s6-rc.d/svc-grackle/run
#!/command/with-contenv bash

exec s6-setuidgid node grackle serve --allow-network
